# scripts/etc/systemd/system/powersave-torbox.service
[Unit]
Description=Powersave script for TORBox

[Service]
ExecStart=/usr/local/bin/powersave.sh

[Install]
WantedBy=default.target