# scripts/etc/init.d/powersave
#! /bin/sh
### BEGIN INIT INFO
# Provides:          powersave
# Required-Start:    $network $syslog
# Required-Stop:     $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: powersave settings
# Description:       sets powersave generated by powertop in /usr/local/bin/powersave.sh
### END INIT INFO

# Do NOT "set -e"
# PATH should only include /usr/* if it runs after the mountnfs.sh script
PATH=/sbin:/usr/sbin:/bin:/usr/bin
DAEMON=/usr/local/bin/powersave.sh
USER=root
STARTSCRIPT=$DAEMON
DESC="powersave"
NAME=powersave
SCRIPTNAME=/etc/init.d/$NAME

test -x $DAEMON || exit 0

# Load the VERBOSE setting and other rcS variables
. /lib/init/vars.sh
# I like to know what is going on
VERBOSE=yes


. /lib/lsb/init-functions

case "$1" in
    start)
        log_daemon_msg "Setting powersave on" "powersave"
        su -c "$STARTSCRIPT" $USER 
        log_end_msg $?
    ;;
  stop)
        log_daemon_msg "Doing nothing" "powersave"
        log_end_msg $?
    ;;
  force-reload|restart)
    $0 stop
    $0 start
    ;;
  status)
    echo "No status, this is oneshot script..."
    ;;
  *)
    echo "Usage: /etc/init.d/$SCRIPTNAME {start|stop|restart|force-reload}"
    exit 1
    ;;
esac

exit 0
